# sasori.service
[Unit]
Description=Sasori Agent Orchestrator
After=network.target

[Service]
Type=simple
User=sasori
Group=sasori
# Directory where the sasori database and handlers/ folder will live
WorkingDirectory=/opt/sasori

# Ensure your PATH includes where sasori was installed
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# ExecStart launches the sasori daemon. Adjust the path if it is not in your global PATH.
ExecStart=/usr/local/bin/sasori

StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=sasori-daemon
Restart=always
RestartSec=10

# Configure your mailbox credentials
Environment="IMAP_SERVER=imap.gmail.com"
Environment="SMTP_SERVER=smtp.gmail.com"
Environment="EMAIL_ACCOUNT=your_agent_email@gmail.com"
Environment="EMAIL_PASSWORD=your_google_app_password"
Environment="POLL_INTERVAL=15"
# Set limits to protect your hardware VRAM
Environment="MAX_CONCURRENT_AGENTS=1"
# Comma-separated list of allowed sender emails
Environment="WHITELIST_EMAILS=your_personal_email@example.com"

[Install]
WantedBy=multi-user.target
